# anycubic_cloud_mqtt_proxy/rootfs/etc/services.d/mqtt-proxy/run
#!/usr/bin/with-contenv sh
# s6 service: inicia o proxy Python e mantém em foreground

echo "[s6] Iniciando Anycubic Cloud MQTT Proxy"
exec python3 /opt/mqtt_proxy/main.py
